# profile.d/brew.sh
#!/usr/bin/env bash
/usr/local/bin/init-brew
eval "$(/home/linuxbrew/.linuxbrew/bin/brew shellenv)"
