# systemd/brew-container.service
[Unit]
Description=Homebrew Distrobox

[Service]
Type=oneshot
ExecStart=/usr/bin/distrobox enter brew-container
RemainAfterExit=true
